# etc/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:/sbin/rungetty tty1 --autologin root
